# {{cookiecutter.repo_name}}/supervisord.conf
[program:{{cookiecutter.repo_name}}]
command=docker-compose up
directory=/code/{{cookiecutter.repo_name}}
redirect_stderr=true
autostart=true
autorestart=true
priority=10